# deploy/systemd/mog-runtime-cleanup.service
[Unit]
Description=Clean stale MOG runtime artifacts

[Service]
Type=oneshot
ExecStart=/srv/mog-template/scripts/cleanup-runtime-artifacts.sh
