# deploy/systemd/mog-runtime-cleanup.timer
[Unit]
Description=Run MOG runtime artifact cleanup daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
